# tools/synchronizer/digiq.timer
[Unit]
Description=DigiQuest synchronizer job timer

[Timer]
OnUnitActiveSec=20s
OnBootSec=20s

[Install]
WantedBy=timers.target